# check-in-bot.service
[Unit]
Description=Python app for automatically adding emoji reacts to slack channels, runs on port 3000
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nginx
Group=nginx
WorkingDirectory=/var/www/check-in-bot
ExecStart=/var/www/check-in-bot/.venv/bin/gunicorn --workers 2 --bind 127.0.0.1:3000 --timeout 120 wsgi:application
Restart=always
RestartSec=10s
StartLimitInterval=60s
StartLimitBurst=5

# Resource limits
LimitNOFILE=65536

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=slack-bot

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
