# check-in-bot-cron.service
[Unit]
Description=Check-in Bot Cron Tasks
After=network.target

[Service]
Type=oneshot
User=www-data
Group=www-data
WorkingDirectory=/var/www/check-in-bot
ExecStart=/var/www/check-in-bot/.venv/bin/python3 cron.py

[Install]
WantedBy=multi-user.target 